# profile.d/osx
# Prepends the given path to PATH environment variable and removes duplicates
prependpath () {
    PATH="$1:${PATH}:"
    PATH="${PATH//:$1}"
    PATH="${PATH%:}"
}

# Puts the latest version of Homebrew install of Emacs into PATH, and makes
# it the default VISUAL editor
emacs_path_guess=/usr/local/Cellar/emacs-plus/*/bin
if [ -x ${emacs_path_guess}/emacsclient ]; then
    prependpath "$emacs_path_guess"
    export VISUAL="emacsclient -c"
    export ALTERNATE_EDITOR=""
    alias emacs="$VISUAL"
fi
